# crontab
0 4 * * * * rm -rf /judge_server/round/*
